# roles/webserver/files/wiki-run.service
[Unit]
Description=wiki-run
[Service]
ExecStart=/usr/local/bin/wiki-run
