# roles/webserver/files/wiki-run.path
[Unit] 
Description=wiki-run path
[Path]
PathModified=/dev/shm/run-wiki
[Install]
WantedBy=multi-user.target
