# roles/webserver/files/www-run.service
[Unit]
Description=www-run
[Service]
ExecStart=/usr/local/bin/www-run
