# roles/webserver/files/www-run.path
[Unit] 
Description=www-run path
[Path]
PathModified=/dev/shm/run-www
[Install]
WantedBy=multi-user.target
