# roles/shell/files/gopher-run.path
[Unit] 
Description=gopher-run path
[Path]
PathModified=/dev/shm/gopher-ansible
[Install]
WantedBy=multi-user.target
