# roles/shell/files/gopher-run.service
[Unit]
Description=gopher-run
[Service]
ExecStart=/usr/local/bin/gopher-run
