# roles/common/files/ansible-run.path
[Unit] 
Description=ansible-run path
[Path]
PathModified=/dev/shm/run-ansible
[Install]
WantedBy=multi-user.target
