# roles/common/files/ansible-run.service
[Unit]
Description=ansible-run
[Service]
ExecStart=/usr/local/bin/ansible-run
