# roles/shell/files/etc/cron.daily/postfix_stats
#!/bin/bash
/usr/sbin/pflogsumm /var/log/mail.log.1 | mail -s "Mail Report for `date +%Y/%m/%d`" postmaster@thunix.net
